# data/supervisord.conf
[supervisord]
nodaemon=true
loglevel=info

[program:Updater]
; Use shell directly to bypass permission issues
command=sh -c "sh /data/entrypoint.sh | awk '/./ {print \"[ UPDATE CRON JOB ]\", $0}'"
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stdout
stderr_logfile_maxbytes=0
redirect_stderr=true

[program:api-dashboard]
; Same pattern as Updater
command=sh -c "pwsh -f /data/www/Start-Server.ps1 | awk '/./ {print \"[ API / Dashboard ]\", $0}'"
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stdout
stderr_logfile_maxbytes=0
redirect_stderr=true

[program:auto-backup]
; Auto backup service that checks for stack changes every 20 sec
command=sh -c "pwsh -f /data/Auto-Backup.ps1 | awk '/./ {print \"[ AUTO BACKUP ]\", $0}'"
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stdout
stderr_logfile_maxbytes=0
redirect_stderr=true
